# ped.socket
[Unit]
Description=ped socket

[Socket]
ListenStream=/run/ped.sock

[Install]
WantedBy=sockets.target